# systemd/bing-wallpaper.timer
[Unit]
Description=Bing Today Wallpaper

[Timer]
OnBootSec=3min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target